# aws/setup-scripts/realtime/curveball-realtime.service
[Unit]
Description=curveball-realtime the realtime server for curveball
After=network.target

[Service]
Environment=REDIS_HOST=cache.dev.curveball.tv
Environment=REDIS_PORT=6379
Environment=NODE_ENV=production
Environment=JWT_SECRET=sEcReT
Environment=INTERNAL_SECRET=VerySecret
Environment=ALLOWED_ORIGINS=foo
Environment=SSL_KEY=/home/ec2-user/realtime/etc/key.pem
Environment=SSL_CERT=/home/ec2-user/realtime/etc/server.crt
Environment=QT_SECRET=CM03sl37vjt2V47eY/EaNWM41wmdWsBxFRqY8/TVpnBEpioDN35oaF2J6bc6M5XQ+7OUdih4PMwCYf6W+/puCMv1upVoFFH22EFbGTSbl/+hJieYt6KAmT9w3OW2xl37BooHQy7yb6crNcZj8hO3pNE/svbpFt9W9tMVpnCmKv00uai2sDmr9fW84TSpQQz4Tavo3pQia+8a5drdjZd3jGcXRft705qa3pVbFtgY5g6cwdpxc7uX0xQdBgHTv/YY7EMI+118RoY6eruLYoCV9BzMKvF+EfmJrl9EGcuNmh+zwE4EKRmX2lghiNLxyGk1ZjFuEBj8CKRkLuWGslKlw==
Type=simple
ExecStart=/usr/bin/npm run start:prod --prefix /home/ec2-user/realtime
Restart=on-failure

[Install]
WantedBy=multi-user.target